# Build-systems/yocto/sysmon/meta-sysmon/recipes-sysmon/sysmon-logger/files/sysmon-logger.service
[Unit]
Description=System Monitor Data Logger
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/sysmon-logger
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
