# Build-systems/yocto/tempmon/meta-tempmon/recipes-tempmon/tempmon/files/tempmon.service
[Unit]
Description=Temperature Monitoring Daemon
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/tempmon
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
