# Build-systems/yocto/sysmon/meta-sysmon/recipes-sysmon/sysmon-web/files/sysmon-web.service
[Unit]
Description=System Monitor Web Dashboard
After=network.target sysmon-logger.service
Requires=sysmon-logger.service

[Service]
Type=simple
ExecStart=/usr/bin/sysmon-web
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
